# scripts/postgres-backup.timer
[Unit]
Description=Run PostgreSQL backup daily at 2:30 AM
Requires=postgres-backup.service

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
